# stage3/01-signalk/files/signalk.socket
[Socket]
ListenStream=443
ListenStream=80

[Install]
WantedBy=sockets.target
