# install-scripts/4-server/files/pypilot_boatimu.service
[Unit]
Description=pypilot boatimu
DefaultDependencies=false
Conflicts=pypilot@pypilot.service

[Service]
Type=simple
ExecStart=pypilot_boatimu -q
StandardOutput=syslog
StandardError=syslog
User=pypilot
Restart=always
RestartSec=2

[Install]
WantedBy=local-fs.target
