# install-scripts/4-server/files/raspotify.socket
[Socket]
ListenStream=41447

[Install]
WantedBy=sockets.target
