# lysmarine/60-firstrun/files/firstRun.service
[Service]
ExecStart=/usr/bin/firstrun
StandardOutput=syslog
StandardError=syslog
User=root

[Install]
WantedBy=graphical.target
